# packages/darwin/Library/LaunchAgents/com.user.homebrew-ssh-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.homebrew-ssh-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ssh-agent</string>
        <string>-D</string>
        <string>-a</string>
        <string>/Users/grim_reaper/.ssh/agent.sock</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SSH_SK_HELPER</key>
        <string>/Users/grim_reaper/.ssh/sk-helper-wrapper</string>
        <key>SSH_ASKPASS</key>
        <string>/Users/grim_reaper/.ssh/sk-askpass</string>
        <key>SSH_ASKPASS_REQUIRE</key>
        <string>force</string>
        <key>HOME</key>
        <string>/Users/grim_reaper</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.user.homebrew-ssh-agent.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.user.homebrew-ssh-agent.err</string>
</dict>
</plist>
